# docker/conf/supervisord.conf
[program:php-fpm]
command=php-fpm
autostart=true
autorestart=true

[program:laravel-queue]
command=php /var/www/html/artisan queue:work --sleep=3 --tries=3
autostart=true
autorestart=true
numprocs=1
